# deploy/systemd/folstingx.service
[Unit]
Description=FolstingX Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/folstingx/backend
ExecStart=/opt/folstingx/bin/folstingx-server
Restart=always
RestartSec=5
User=folstingx
Group=folstingx
LimitNOFILE=65535
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
